# app/fscn8001.service
[Unit]
Description=赋安 FSCN8001 用传协议解析服务（平台集成版）
Documentation=https://github.com/xzy-fire-platform
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fsan
ExecStart=/usr/bin/node /opt/fsan/fscn8001.js
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=append:/var/log/fscn8001.log
StandardError=append:/var/log/fscn8001.log

# 环境变量
Environment="NODE_ENV=production"
Environment="PUSH_URL=http://127.0.0.1:5003/api/fscn8001/push"

# 资源限制
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
